# etc/autosnap.timer
[Unit]
Description=Automatically snapshot ZFS datasets every 15 minutes

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
