# etc/autosnap.service
[Unit]
Description=Automatically snapshot ZFS datasets
Requires=zfs.target
After=zfs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/autosnap -v
